# Package_all/DEBIAN/postinst
systemctl enable pi-deck
systemctl start pi-deck

