# Package_all/DEBIAN/prerm
systemctl stop pi-deck
systemctl disable pi-deck
